# files/supervisor_worker.conf.j2
[program:pycsw-load]
command=/home/pycsw/current/.venv/bin/pycsw-load.sh
stdout_logfile=/var/log/pycsw/pycsw-load.log
redirect_stderr=true
autostart=false
autorestart=false
